# systemd/nutrition-agent.service
[Unit]
Description=nutrition-agent Telegram bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=nutrition-agent
WorkingDirectory=/opt/nutrition-agent
EnvironmentFile=/opt/nutrition-agent/.env
ExecStart=/opt/nutrition-agent/.venv/bin/python -m app.bot.telegram_bot
Restart=on-failure
RestartSec=10
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
